# systemd/fsbackup-s3-export.timer
[Unit]
Description=Run fsbackup S3 export nightly

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
